# mysql-5.1.73_amd64/mysql-legacy.conf
# mysql-legacy.conf
start on filesystem
exec /usr/local/mysql/bin/mysqld_safe
